# ansible/roles/server/templates/pixel-perfect-web.service.j2
[Unit]
Description=PixelPerfect Web Server

[Service]
WorkingDirectory={{ document_root }}/backend
ExecStart=/usr/bin/npm start
Restart=on-abort

[Install]
WantedBy=multi-user.target